# wa_gateway/wa-gateway.service
[Unit]
Description=WhatsApp Gateway — whatsapp-web.js bridge for LeadQualBot
After=network.target agent.service

[Service]
Type=simple
User=selfops
WorkingDirectory=/personalized_agent_bot/wa_gateway
ExecStart=/usr/bin/node gateway.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
